# monitor/topLogger.service
[Unit]
Description=logs_data_from_game_servers

[Service]
ExecStart=/home/ark/epicsrvrmgr/monitor/logArkResourceUsage.sh
